# examples/systemd/temporal-worker.service
[Unit]
Description=ActiveJob Temporal Worker
After=syslog.target network.target

[Service]
Type=simple
WorkingDirectory=/opt/myapp/current

# Adjust based on your Ruby setup:
# For rbenv:
# ExecStart=/home/deploy/.rbenv/shims/bundle exec bin/temporal-worker
# For rvm:
# ExecStart=/home/deploy/.rvm/bin/rvm in /opt/myapp/current do bundle exec bin/temporal-worker
# For system ruby:
ExecStart=/usr/local/bin/bundle exec bin/temporal-worker

# Environment Configuration
Environment="TEMPORAL_TARGET=temporal.example.com:7233"
Environment="TEMPORAL_NAMESPACE=production"
Environment="AJ_TEMPORAL_WORKER_QUEUE=default"
Environment="RAILS_ROOT=/opt/myapp/current"
Environment="RAILS_ENV=production"

# Concurrency and Performance
Environment="AJ_TEMPORAL_MAX_ACT=100"

# Service User
User=deploy
Group=deploy

# Restart Policy
Restart=always
RestartSec=5s

# Timeouts
StartLimitInterval=600
StartLimitBurst=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=temporal-worker

# Resource Limits (adjust as needed)
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
